# online-api/supervisord/cherrypy.conf
[supervisord]

[program:skillsapi]
command=gunicorn wsgi --config /opt/skills-online-api/gunicorn.conf --chdir=/opt/skills-online-api
